# scripts/udev/60-monitor-plug.rules
# Update xrandr configuration
ACTION=="change", SUBSYSTEM=="drm", RUN+="/bin/bash /home/vitdevelop/.config/scripts/udev/calibrate_monitors.sh"
